# pbxboot/includes.chroot/etc/systemd/system/getty@tty1.service.d/console.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --login-program /usr/local/bin/phonebocx.sh --noclear %I $TERM
TTYVTDisallocate=no

